# device-support/arduino-portenta-x8/yocto/meta-normacore-x8/recipes-connectivity/tailscaled-init/files/tailscaled
#!/bin/sh

### BEGIN INIT INFO
# Provides:          tailscaled
# Required-Start:    $network $remote_fs
# Required-Stop:     $network $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Tailscale daemon
### END INIT INFO

PIDFILE=/run/tailscaled.pid
STATE=/var/lib/tailscale/tailscaled.state
SOCKET=/run/tailscale/tailscaled.sock

find_daemon() {
    for p in /usr/sbin/tailscaled /usr/bin/tailscaled /sbin/tailscaled /bin/tailscaled; do
        [ -x "$p" ] && echo "$p" && return 0
    done
    command -v tailscaled 2>/dev/null || true
}

DAEMON="$(find_daemon)"

case "$1" in
    start)
        [ -n "$DAEMON" ] || {
            echo "tailscaled not found"
            exit 1
        }

        mkdir -p /run/tailscale /var/lib/tailscale /dev/net
        [ -c /dev/net/tun ] || mknod /dev/net/tun c 10 200 2>/dev/null || true
        modprobe tun 2>/dev/null || true

        echo "Starting tailscaled"
        start-stop-daemon -S -b -m -p "$PIDFILE" -x "$DAEMON" -- \
            --state="$STATE" \
            --socket="$SOCKET"
        ;;

    stop)
        echo "Stopping tailscaled"
        start-stop-daemon -K -p "$PIDFILE" 2>/dev/null || true
        rm -f "$PIDFILE"
        ;;

    restart)
        "$0" stop
        sleep 1
        "$0" start
        ;;

    status)
        if [ -f "$PIDFILE" ] && kill -0 "$(cat "$PIDFILE")" 2>/dev/null; then
            echo "tailscaled running"
            exit 0
        fi
        echo "tailscaled not running"
        exit 3
        ;;

    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
